# Web/xssend_me_flag/cronjob/crontab
PATH="/usr/bin/"

*/1 * * * * node /app/puppeteer.js
